# snapserver/rootfs/etc/s6-overlay/s6-rc.d/snapserver/run
#!/command/with-contenv bash

exec /run.sh
