# snapserver/rootfs/etc/services.d/disable-dmsg/run
#!/command/with-contenv bashio

bashio::log.info "[dmesg] Reducing kernel console verbosity"
exec dmesg -n 1
