# snapserver/rootfs/etc/services.d/disable-dmsg/run
#!/command/with-contenv bashio

bashio::log.info "[dmesg] Reducing kernel console verbosity"
if ! dmesg -n 1; then
    bashio::log.warning "[dmesg] Unable to adjust kernel console verbosity"
fi

exit 0
