# snapserver/rootfs/etc/services.d/disable-dmsg/run
#!/command/with-contenv bash
# Disable kernel logs to console
exec dmesg --console-off
