# snapserver/rootfs/etc/services.d/pulseaudio/run
#!/usr/bin/execlineb -P
with-contenv
pulseaudio --exit-idle-time=-1 --system --disallow-exit --log-target=stderr --daemonize=no