# snapserver/rootfs/etc/services.d/pulseaudio/run
#!/command/execlineb -P
/command/with-contenv
pulseaudio --exit-idle-time=-1 --system --disallow-exit --log-target=stderr --daemonize=no -nF /etc/pulse/system.pa

