# snapserver/rootfs/etc/s6-overlay/s6-rc.d/pulseaudio/run
#!/command/with-contenv bash

set -euo pipefail

log() {
    local level="$1"
    shift
    printf '[%s] %s\n' "$level" "$*"
}

log_info() {
    log "INFO" "$@"
}

log_warning() {
    log "WARN" "$@"
}

log_info "[PA] Preparing runtime directories"
mkdir -p /var/run/pulse/.config/pulse
mkdir -p /run/pulse
mkdir -p /etc/pulse/system.pa.d

if [[ ! -f /var/run/pulse/.config/pulse/cookie ]]; then
    head -c 256 /dev/urandom > /var/run/pulse/.config/pulse/cookie || true
fi

ln -sf /var/run/pulse/.config/pulse/cookie /var/run/pulse/.pulse-cookie

chmod 600 /var/run/pulse/.config/pulse/cookie /var/run/pulse/.pulse-cookie 2>/dev/null || true

chown -R pulse:pulse /var/run/pulse /run/pulse || true

export PULSE_RUNTIME_PATH="/var/run/pulse"
export XDG_RUNTIME_DIR="/var/run/pulse"
export HOME="/var/run/pulse"
export XDG_CONFIG_HOME="/var/run/pulse/.config"
export PULSE_COOKIE="/var/run/pulse/.config/pulse/cookie"
export DBUS_SYSTEM_BUS_ADDRESS="unix:path=/var/run/dbus/system_bus_socket"
export DBUS_SESSION_BUS_ADDRESS="${DBUS_SYSTEM_BUS_ADDRESS}"

log_info "[PA] Waiting for system bus"
for _ in $(seq 1 50); do
    if [[ -S /var/run/dbus/system_bus_socket ]]; then
        break
    fi
    sleep 0.1
done

if [[ ! -S /var/run/dbus/system_bus_socket ]]; then
    log_warning "[PA] system bus socket not available after waiting"
fi

log_info "[PA] Starting PulseAudio"
if command -v setpriv >/dev/null 2>&1; then
    exec setpriv --reuid pulse --regid pulse --init-groups \
        pulseaudio --exit-idle-time=-1 --disallow-exit --disallow-module-loading \
        --log-target=stderr --daemonize=no -nF /etc/pulse/system.pa
fi

if command -v runuser >/dev/null 2>&1; then
    exec runuser -u pulse -- pulseaudio --exit-idle-time=-1 --disallow-exit --disallow-module-loading \
        --log-target=stderr --daemonize=no -nF /etc/pulse/system.pa
fi

if command -v s6-setuidgid >/dev/null 2>&1; then
    exec s6-setuidgid pulse pulseaudio --exit-idle-time=-1 --disallow-exit --disallow-module-loading \
        --log-target=stderr --daemonize=no -nF /etc/pulse/system.pa
fi

# As a last resort, run PulseAudio as root so the service can still start.
exec pulseaudio --exit-idle-time=-1 --disallow-exit --disallow-module-loading \
    --log-target=stderr --daemonize=no -nF /etc/pulse/system.pa
