# snapserver/rootfs/etc/services.d/bluetooth/run
#!/command/with-contenv bashio

bashio::log.info "[BT] Starting bluetoothd"
exec /usr/lib/bluetooth/bluetoothd --noplugin=sap --nodetach
