# snapserver/rootfs/etc/services.d/bluetooth/run
#!/usr/bin/with-contenv bash
echo "[BT] Starting bluetoothd" > /proc/1/fd/1

# Starta dbus om den inte redan kör
if [ ! -e /var/run/dbus ]; then
  mkdir -p /var/run/dbus
fi
dbus-daemon --system --fork

# Starta bluetoothd i bakgrunden
/usr/lib/bluetooth/bluetoothd --noplugin=sap &

# Vänta lite så att bluetoothd hinner starta
sleep 2

# Leta efter första tillgängliga controller
controller=$(bluetoothctl list | head -n1 | awk '{print $2}')

if [ -z "$controller" ]; then
  echo "[BT] Ingen Bluetooth-controller hittades ännu"
else
  echo "[BT] Hittade controller: $controller"
  echo -e "select $controller\npower on\nagent on\ndefault-agent\npairable on\ndiscoverable on" | bluetoothctl
fi

# Håll processen igång
tail -f /dev/null