# snapserver/rootfs/etc/services.d/bluetooth/run
#!/command/with-contenv bashio

bashio::log.info "[BT] Waiting for system bus"
for _ in $(seq 1 50); do
    if [[ -S /var/run/dbus/system_bus_socket ]]; then
        break
    fi
    sleep 0.1
done

if [[ ! -S /var/run/dbus/system_bus_socket ]]; then
    bashio::log.warning "[BT] system bus socket not available after waiting"
fi

bashio::log.info "[BT] Starting bluetoothd"
if pidof bluetoothd >/dev/null 2>&1; then
    bashio::log.warning "[BT] Detected an existing bluetoothd process on the host. Skipping internal daemon startup."
    exec tail -f /dev/null
fi

exec /usr/lib/bluetooth/bluetoothd --noplugin=sap --nodetach
