# snapserver/rootfs/etc/services.d/bluetooth/run
#!/usr/bin/with-contenv bash
echo "[BT] Starting bluetoothd" > /proc/1/fd/1

# Starta dbus om den inte redan kör
if [ ! -e /var/run/dbus ]; then
  mkdir -p /var/run/dbus
fi

# Starta dbus i bakgrunden (utan fork)
dbus-daemon --system --nofork --nopidfile --nosyslog &
sleep 1

# Starta bluetoothd i bakgrunden
/usr/lib/bluetooth/bluetoothd --noplugin=sap &

# Vänta lite så att bluetoothd hinner starta
sleep 2

# Skapa fifo för Snapcast om den inte redan finns
if [ ! -p /tmp/snapfifo ]; then
  mkfifo /tmp/snapfifo
  echo "[BT] Created /tmp/snapfifo"
fi

# Starta pulseaudio i bakgrunden (utan -D)
pulseaudio --exit-idle-time=-1 --system --disallow-exit --log-target=stderr &
sleep 2

# Skapa en virtual sink (pipe) om den inte redan finns
if ! pactl list short sinks | grep -q bt2snap; then
  pactl load-module module-pipe-sink file=/tmp/snapfifo format=s16le rate=44100 channels=2 sink_name=bt2snap > /dev/null
fi

# Sätt bt2snap som default-sink så att Bluetooth-ljudet hamnar rätt
pactl set-default-sink bt2snap

# Leta efter första tillgängliga controller
controller=$(bluetoothctl list | head -n1 | awk '{print $2}')

if [ -z "$controller" ]; then
  echo "[BT] Ingen Bluetooth-controller hittades ännu"
else
  echo "[BT] Hittade controller: $controller"
  echo -e "select $controller\npower on\nagent on\ndefault-agent\npairable on\ndiscoverable on" | bluetoothctl
fi

# Håll processen igång
exec tail -f /dev/null