# snapserver/rootfs/etc/services.d/bluetooth/run
#!/command/with-contenv bash
echo "[BT] Starting bluetoothd" > /proc/1/fd/1
exec /usr/lib/bluetooth/bluetoothd --noplugin=sap --nodetach
