# snapserver/rootfs/etc/s6-overlay/s6-rc.d/avahi/run
#!/command/with-contenv bash

set -euo pipefail

log() {
    local level="$1"
    shift
    printf '[%s] %s\n' "$level" "$*"
}

log_info() {
    log "INFO" "$@"
}

log_warning() {
    log "WARN" "$@"
}

log_info "[Avahi] Waiting for system bus"
for _ in $(seq 1 50); do
    if [[ -S /var/run/dbus/system_bus_socket ]]; then
        break
    fi
    sleep 0.1
done

if [[ ! -S /var/run/dbus/system_bus_socket ]]; then
    log_warning "[Avahi] system bus socket not available after waiting"
fi

log_info "[Avahi] Starting avahi-daemon"
exec avahi-daemon --no-chroot
