# snapserver/rootfs/etc/s6-overlay/s6-rc.d/dbus/run
#!/command/with-contenv bash

set -euo pipefail

log() {
    local level="$1"
    shift
    printf '[%s] %s\n' "$level" "$*"
}

log_info() {
    log "INFO" "$@"
}

log_warning() {
    log "WARN" "$@"
}

log_info "[DBUS] Preparing runtime directory"
mkdir -p /run/dbus

if [[ ! -e /var/run/dbus ]]; then
    mkdir -p /var/run
    ln -s /run/dbus /var/run/dbus
fi

log_info "[DBUS] Starting system bus"
exec dbus-daemon --system --nofork --nopidfile --nosyslog
