# snapserver/rootfs/etc/services.d/snapserver/run
#!/command/with-contenv bashio

exec /run.sh
