# snapserver/rootfs/etc/s6-overlay/s6-rc.d/legacy-services/run
#!/command/with-contenv bash
set -euo pipefail

# Ensure the legacy services directory exists for compatibility.
if [[ ! -d /etc/services.d ]]; then
    mkdir -p /etc/services.d
fi

# Run an s6 supervisor for legacy service definitions without relying on
# s6-overlay-suexec, which now requires PID 1.
exec /command/s6-svscan -t0 /etc/services.d
