# snapserver/rootfs/etc/s6-overlay/s6-rc.d/legacy-services/up
#!/command/with-contenv sh

printf '%s\n' "[legacy] Skipping deprecated legacy service supervision"
