# snapserver/rootfs/etc/s6-overlay/s6-rc.d/legacy-services/up
#!/command/with-contenv sh

# The services within this add-on are implemented with native s6-rc entries.
# Keep the legacy shim in a no-op state so the supervisor doesn't attempt to
# launch the deprecated services.d tree.
printf '%s\n' "[INFO] [s6] legacy services shim disabled"
exit 0
