# snapserver/rootfs/etc/services.d/dbus/run
#!/usr/bin/execlineb -P
with-contenv
dbus-daemon --system --nofork --nopidfile --nosyslog