# snapserver/rootfs/etc/services.d/dbus/run
#!/command/execlineb -P
/command/with-contenv
dbus-daemon --system --nofork --nopidfile --nosyslog

