# snapserver/rootfs/etc/s6-overlay/s6-rc.d/disable-dmesg/up
#!/command/with-contenv sh

printf '%s\n' "[dmesg] Reducing kernel console verbosity"
if ! dmesg -n 1 >/dev/null 2>&1; then
    printf '%s\n' "[dmesg] Unable to adjust kernel console verbosity"
fi
